# config/laravel-boot.service
[Unit]
Description=Laravel boot tasks (migrations)
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=exedev
ExecStart=/usr/local/bin/laravel-boot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
